# zookeeper/zookeeper.service
[Unit]
Description=miner

[Service]
ExecStart=/path/to/miner
ExecStop=
ExecReload=

[Install]
WantedBy=multi-user.target
